# com.tavr-digest.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tavr-digest.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mbowdish/tavr-digest/venv/bin/python3</string>
        <string>/Users/mbowdish/tavr-digest/main.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mbowdish/tavr-digest</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/mbowdish/tavr-digest/data/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mbowdish/tavr-digest/data/launchd_stderr.log</string>
</dict>
</plist>
